# cronjobs/crontab
SHELL=/bin/bash
PYTHONPATH=/usr/local/bin/://usr/local/lib/python3.7/site-packages
0 0 * * * python3 /code/manage.py fetchdata >/proc/1/fd/1 2>/proc/1/fd/2
